# configuration/systemd/services/check-linux-common.timer
[Unit]
Description=check linux-common git repo for updates

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
